# set-coall.service
[Unit]
Description=Override AMD Curve Optimizer

[Service]
Type=simple
ExecStart=/usr/local/bin/ryzenadj --set-coall=0x100014

